# deploy/s6-overlay/s6-rc.d/queue-worker/run
#!/command/with-contenv bash
#TODO: Change to work
s6-setuidgid webuser php $WEBUSER_HOME/artisan queue:listen --tries=3
